# kernel/udev/test.service
[Unit]
Description="for testing"

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/home/linux/test.s

[Install]
WantedBy=multi-user.target
